# Dockerfile/supervisord.conf
[program:sshd]
command=/usr/sbin/sshd -D
autostart=true
autorestart=true

[program:docker]
command=docker -d -dns 172.17.42.1 -H unix:///var/run/docker.sock -r=false
autostart=true
autorestart=true
